# supervisor_505aq.ini
[program:505aq]
command = /home/jiakai/505a-air-quality/daemon.py
redirect_stderr = true
